# Yocto-image/meta-apps/recipes-osm/osm/osm/db.service
[Unit]
Description=OSM Database Handler Service
After=network.target
Wants=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/usr/bin/OSM
ExecStart=/usr/bin/python3 -u /usr/bin/OSM/initialize_db.py
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0  

[Install]
WantedBy=multi-user.target
